# deploy/secretary-gateway.service
[Unit]
Description=Secretary AI Gateway (Reverse Proxy)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/john/projects/secretary
ExecStart=/usr/bin/python3 -m scripts.ai_gateway
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secretary-gateway
EnvironmentFile=/home/john/projects/secretary/bot/.env
NoNewPrivileges=yes

[Install]
WantedBy=default.target
