# deploy/secretary-listener.service
[Unit]
Description=Secretary Telegram Listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=john
Group=john
WorkingDirectory=/home/john/projects/secretary
ExecStart=/usr/bin/python3 -m bot.telegram_listener
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secretary-listener
EnvironmentFile=/home/john/projects/secretary/bot/.env

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/john/projects/secretary/bot
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
