# deploy/secretary-worker.service
[Unit]
Description=Secretary Message Worker
After=network-online.target secretary-listener.service
Wants=network-online.target

[Service]
Type=simple
User=john
Group=john
WorkingDirectory=/home/john/projects/secretary
ExecStart=/usr/bin/python3 -m bot.worker
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secretary-worker
EnvironmentFile=/home/john/projects/secretary/bot/.env

# Worker needs claude CLI + node in PATH
Environment=PATH=/home/john/.local/bin:/home/john/.nvm/versions/node/v20.19.6/bin:/usr/local/bin:/usr/bin:/bin:/home/john/.local/share/claude/versions
# Prevent "cannot launch inside another Claude Code session" error
Environment=CLAUDE_CODE_ENTRY_POINT=
Environment=CLAUDECODE=

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/john/projects/secretary/bot
ReadWritePaths=/home/john/.claude
ReadWritePaths=/home/john/.local/share/claude
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
